# supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock          ; Matches the original image's socket path
chmod=0700                             ; Keep it secure (owner-only)
username=admin                         ; ← Add this (change to whatever you want)
password=change_me_strong_password     ; ← Add this (make it secure!)

; The rest can stay minimal/default or copy more from the original if needed
[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true                          ; Important for Docker
minfds=1024
minprocs=200
user=root                              ; Usually root in this image

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock
username=admin                         ; Must match above
password=admin987654321     ; Must match above

; Keep the programs (nginx + php-fpm) as in the original
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
